# src/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn --config=gunicorn_conf.py api:app -b 0.0.0.0:8080
directory=/ailab/src
autostart=true
autorestart=true
redirect_stderr=True